# .config/systemd/user/clock.service
[Unit]
Description=Notify 10 minutes till
[Service]
ExecStart=bash -c "~/bin/notify-time"
