# .config/systemd/user/updates.service
[Unit]
Description=Notify when updates are available
[Service]
ExecStart=bash -c "~/bin/upgrade-notify"
