# .config/systemd/user/updates.timer
[Unit]
Description=Check for updates
[Timer]
OnBootSec=5 min
OnUnitActiveSec=30 min
[Install]
WantedBy=multi-user.target
