# udev/95-monitors.rules
ACTION=="change",RUN+="/home/ktravis/bin/displays"
